# deploy/mcp-ssh.service
[Unit]
Description=mcp-ssh MCP server (executes shell commands locally)
After=network.target

[Service]
Type=simple
User=mcp-ssh
Group=mcp-ssh
ExecStart=/usr/bin/mcp-ssh serve
# Writes job dir `/var/lib/mcp-ssh/logs/jobs`. StateDirectory creates the dir
# writable by the service user even with ProtectSystem=full.
StateDirectory=mcp-ssh
# Optional env overrides. Config credentials live at `/etc/mcp-ssh/config.toml`
# (written by `mcp-ssh set-auth`); this file is for ephemeral knobs.
EnvironmentFile=-/etc/mcp-ssh/mcp-ssh.env
Restart=on-failure
RestartSec=5

# Hardening. NOTE: this tool intentionally shells out to run arbitrary
# commands as the service user, so sandboxing is necessarily loose -- these
# only raise the floor, they do not contain a determined command.
#
# NoNewPrivileges is OFF by design: the agent is meant to self-manage the host
# (install its own updates, restart itself, manage services) via `sudo`, which
# is setuid and blocked by NoNewPrivileges=true. The service user is granted
# sudo by the installer. To lock the box down instead, set this back to true and
# remove /etc/sudoers.d/mcp-ssh -- see docs/deploy.md.
NoNewPrivileges=false
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
